# packages/rpnd-mod-rt433/root/etc/uci-defaults/50-rpnd-mod-rt433
#!/bin/sh

if [ "$(uci -q get rpnd.@rt433[0])" = '' ] ; then
	uci add rpnd rt433
	uci commit
fi

exit 0
