# packages/rpnd-mod-rt433/root/etc/uci-defaults/50_rpnd_mod_rt433
#!/bin/sh

uci -q batch << EOF
add rpnd rt433
commit
exit
EOF
exit 0
