# packages/rpnd-mod-gpio/root/etc/uci-defaults/50-rpnd-mod-gpio
#!/bin/sh

if [ "$(uci -q get rpnd.@gpio[0])" = '' ] ; then
	uci add rpnd gpio
	uci set rpnd.@gpio[0].rpio=0
	uci commit
fi

exit 0
