# packages/rpnd-base/root/etc/uci-defaults/50-rpnd-base
#!/bin/sh

touch /etc/config/rpnd

if [ "$(uci -q get rpnd.@rpnd[0])" = '' ] ; then
 	uci add rpnd rpnd
	uci set wireless.radio0.disabled=0
	uci set wireless.default_radio0.ssid=Rpnd-Setup
	uci set network.lan.ipaddr=192.168.77.1
	uci commit
fi
if [ "$(grep -q Rpnd-initialized "/boot/config.txt")" = '' ] ; then
	printf "dtoverlay=gpio-no-irq\n" >> /boot/config.txt
	printf "#Rpnd-initialized\n" >> /boot/config.txt
fi
exit 0
