# packages/rpnd/root/etc/uci-defaults/50_rpnd
#!/bin/sh
if ! grep -q Rpnd-initialized "/boot/config.txt"; then
	uci -q batch << EOF
	set wireless.radio0.disabled=0
	set wireless.default_radio0.ssid=Rpnd-Setup
	set network.lan.ipaddr=192.168.77.1
	commit
	exit
EOF
	printf "ttyGS0::respawn:/usr/sbin/agetty -h -t 60 ttyGS0 115200 vt102\n" >> /etc/inittab
	printf "dtoverlay=gpio-no-irq\n" >> /boot/config.txt
	printf "dtoverlay=dwc2,dr_mode=otg\n" >> /boot/config.txt
	printf "#Rpnd-initialized\n" >> /boot/config.txt
fi
exit 0
