# packages/rpnd/root/etc/uci-defaults/50_rpnd
#!/bin/sh
uci -q batch << EOF



exit
EOF
exit 0
