# packages/rpnd/root/etc/uci-defaults/50_rpnd
#!/bin/sh
uci -q batch << EOF
set wireless.radio0.disabled=0
set wireless.radio0.ssid=Rpnd-Setup
commit
exit
EOF
echo dtoverlay=gpio-no-irq >> /boot/config.txt
service network restart
exit 0
