# packages/rpnd-mod-wire/root/etc/uci-defaults/50-rpnd-mod-wire
#!/bin/sh

if [ "$(uci -q get rpnd.@wire[0])" = '' ] ; then
	uci add rpnd wire
	uci commit
fi

exit 0
