# packages/rpnd-mod-wire/root/etc/uci-defaults/50_rpnd_mod_wire
#!/bin/sh

uci -q batch << EOF
add rpnd wire
commit
exit
EOF
exit 0
