# packages/rpnd-mod-chime/root/etc/uci-defaults/50-rpnd-mod-chime
#!/bin/sh

if [ "$(uci -q get rpnd.@chime[0])" = '' ] ; then
	uci add rpnd chime
	uci commit
fi

exit 0
