# packages/rpnd-mod-chime/root/etc/uci-defaults/50_rpnd_mod_chime
#!/bin/sh

uci -q batch << EOF
add rpnd chime
commit
exit
EOF
exit 0
