# packages/rpnd-mod-i2c/root/etc/uci-defaults/50-rpnd-mod-ic2
#!/bin/sh

if [ "$(uci -q get rpnd.@i2c[0])" = '' ] ; then
	uci add rpnd i2c
	uci set rpnd.@i2c[0].disabled=0
	uci commit
fi

exit 0
