# packages/rpnd-mod-i2c/root/etc/uci-defaults/50_rpnd_mod_ic2
#!/bin/sh

uci -q batch << EOF
add rpnd i2c
commit
exit
EOF
exit 0
